# pkgbuild/uresourced.install
post_install() {
	systemctl enable --now uresourced
}

post_upgrade() {
	post_install
	systemctl restart uresourced
}
